# install/mandrel-watcher.service
[Unit]
Description=Mandrel Watcher - Git commit sync daemon
Documentation=https://github.com/RidgetopAi/mandrel-watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/mandrel-watcher start -f
ExecStop=/usr/local/bin/mandrel-watcher stop
Restart=on-failure
RestartSec=10

# Run as user (not root)
User=%i

# Environment
Environment=HOME=/home/%i

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mandrel-watcher

[Install]
WantedBy=default.target
